# npm-registry/conf/app.ini
[program:app]
command = npm-register start --always-https
directory = /npm-registry
user = root
autostart = true
autorestart = true
stdout_logfile = /var/log/supervisor/app.log
stderr_logfile = /var/log/supervisor/app_errors.log
environment = NODE_ENV="production" NPM_REGISTER_FS_DIRECTORY=/npm-register/cache