# cloud9/conf/etc/profile.d/ghost.sh
#!/bin/bash

if [ -d /ghost/.bin ]; then
    PATH="/ghost/.bin:$PATH"
fi
