# nginx-static/conf/etc/supervisor/conf.d/nginx.conf
[program:sshd]
command=/usr/sbin/nginx -g 
autorestart=true
startretries=3