# npm-register/conf/app.ini
[program:app]
command = node node_modules/npm-register/lib/start.js
directory = /npm-registry
user = root
autostart = true
autorestart = true
stdout_logfile = /tmp/app.log
stderr_logfile = /tmp/app_errors.log
environment = NODE_ENV="production"