# npm-register/conf/app.ini
[program:app]
command =  node node_modules/npm-register/lib/start.js --always-https
directory = /npm-register
user = root
autostart = true
autorestart = true
stdout_logfile = /var/log/supervisor/app.log
stderr_logfile = /var/log/supervisor/app_errors.log
environment = NODE_ENV="production", NPM_REGISTER_STORAGE=fs, NPM_REGISTER_FS_DIRECTORY="/npm-register/cache"